# cookbooks/zookeeper-component/templates/default/exhibitor.upstart.conf.erb
#!upstart
description "Exhibitor Server"

start on startup
stop on shutdown

respawn

exec su -c "java -jar <%= @jar %> \
    <%= @opts.keys.sort.map { |k| "--#{k.to_s} #{@opts[k]}" }.join(" ") %>" <%= @user %>

post-start script
  <%= @check_script %> || { stop ; exit 1; }
end script
